# supervisor.conf
[program:pinger]
command=/home/justin/code/pinger/venv/bin/python /home/justin/code/pinger/pinger.py
user=justin

